# scripts/fetch.cron
10 * * * * /home/web/CRAWL/bin/python /home/web/crawl/scripts/fetch_latest.py >> /home/web/crawl/scripts/debug.log
